# docker/mbstats.crontab
SHELL=/bin/bash
VIRTUAL_ENV=/opt/venv
PATH=/opt/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""

* * * * * root /mbstats/run -c /etc/container_environment.json --log-handler syslog --workdir /data

