# files/mbstats.crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * root /tools/stats.parser.py -c /etc/container_environment.json --syslog --workdir /data
